# protocols/P25/gateway/p25-gateway.service
# /etc/systemd/system/p25-gateway.service
[Unit]
Description=P25 Gateway
After=p25-host.service
Requires=p25-host.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/OpenMT/protocols/P25
ExecStart=/home/pi/OpenMT/protocols/P25/build/p25_gateway
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
